# ops/launchd/com.stockflow.daily-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.stockflow.daily-refresh</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/saadmanzoor/Documents/StockFlow/scripts/daily-postgres-refresh.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/saadmanzoor/Documents/StockFlow</string>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>6</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/saadmanzoor/Documents/StockFlow/logs/launchd-daily-refresh.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/saadmanzoor/Documents/StockFlow/logs/launchd-daily-refresh.err.log</string>

    <key>RunAtLoad</key>
    <false/>
  </dict>
</plist>
